# ghp_mm_sniffer/rootfs/etc/services.d/ghp-mm-sniffer/run
#!/usr/bin/with-contenv bash

# Активируем виртуальное окружение
source /usr/src/app/venv/bin/activate

# Запуск Python-скрипта
exec python3 /usr/src/app/ghp-mm2mqtt.py
